# meta-yandex/meta-openrack/meta-shaosi/recipes-phosphor/tobby/files/obmc-tobby.service
[Unit]
Description=Tobby daemon
Requires=obmc-shaosid.service openresty.service
[Service]
Type=simple
Restart=always
RestartSec=0
ExecStart=/usr/sbin/tobby
[Install]
WantedBy=network.target
